# etc/restart-controller.service
[Unit]
Description=Restart temperature controller

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart controller.service

